# rel/crontab
51 23 * * 1,4 sh /app/dl-and-import-geolite2.sh
